# deploy/alpha-sidekiq/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:sidekiq]
user=www-data
directory=/app
command=bundle exec sidekiq
